# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:fastapi]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 8080 --log-level info
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:openclaw-init]
; Writes models.json + default model before the gateway starts
command=/app/openclaw-setup.sh
directory=/app
autostart=true
autorestart=false
priority=12
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:openclaw]
; Starts before CM handler — gateway must be up before Python tries to connect
command=node /app/openclaw/openclaw.mjs gateway --port 18789 --allow-unconfigured
directory=/app/openclaw
autostart=true
autorestart=true
priority=15
startsecs=3
environment=OPENAI_BASE_URL="http://localhost:8080/llm/proxy/openai/v1",openai_API_KEY="foundry-proxied",OPENCLAW_GATEWAY_TOKEN="%(ENV_OPENCLAW_GATEWAY_TOKEN)s",OPENCLAW_STATE_DIR="/tmp/openclaw-state",OPENCLAW_SKIP_CHANNELS="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:compute_module]
; Starts last — retries internally until FastAPI is healthy
command=python compute_module/handler.py
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
